# host/udev.rules
# A udev rule to run a custom script when ZFS zvols are added
KERNEL=="zd*" RUN+="/usr/local/bin/zfs-mknod"
